# deploy/systemd/btc-vol-auth.service
[Unit]
Description=BTCAlert dashboard cookie authentication gateway
Wants=network-online.target
After=network-online.target
Before=nginx.service

[Service]
Type=simple
User=btcalert
Group=nginx
SupplementaryGroups=btcalert
WorkingDirectory=/opt/btc-vol-alert
EnvironmentFile=/etc/btc-vol-alert.env
Environment=HOME=/var/lib/btc-vol-alert
RuntimeDirectory=btc-vol-alert
RuntimeDirectoryMode=0750
UMask=0007
ExecStart=/opt/btc-vol-alert/.venv/bin/gunicorn --workers 1 --threads 2 --timeout 30 --bind unix:/run/btc-vol-alert/auth.sock --umask 007 --access-logfile - --error-logfile - "auth_gateway:create_app()"
Restart=on-failure
RestartSec=5s
TimeoutStartSec=30s
TimeoutStopSec=20s
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=true
ReadOnlyPaths=/opt/btc-vol-alert /etc/btc-vol-alert.env /etc/btc-vol-alert-dashboard.htpasswd
LockPersonality=true
RestrictSUIDSGID=true
RestrictRealtime=true

[Install]
WantedBy=multi-user.target
