# deploy/systemd/btc-vol-dashboard.service
[Unit]
Description=BTCAlert read-only Streamlit dashboard
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=btcalert
Group=btcalert
WorkingDirectory=/opt/btc-vol-alert
EnvironmentFile=/etc/btc-vol-alert.env
Environment=HOME=/var/lib/btc-vol-alert
Environment=STATUS_STATE_FILE=/var/lib/btc-vol-alert/state.json
ExecStart=/opt/btc-vol-alert/.venv/bin/streamlit run /opt/btc-vol-alert/status_app.py --server.address 127.0.0.1 --server.port 8503 --server.headless true --browser.gatherUsageStats false
Restart=on-failure
RestartSec=5s
TimeoutStopSec=20s
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadOnlyPaths=/opt/btc-vol-alert
ReadWritePaths=/var/lib/btc-vol-alert /tmp
LockPersonality=true

[Install]
WantedBy=multi-user.target
